# deploy/intraday.service
[Unit]
Description=Recostock Intraday ETF Signal Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/recostock
EnvironmentFile=/root/recostock/deploy/.env
ExecStart=/root/recostock/venv/bin/python scripts/run_intraday.py
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=intraday-bot

[Install]
WantedBy=multi-user.target
