# deploy/intraday.service
[Unit]
Description=Recostock Intraday ETF Signal Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/recostock
EnvironmentFile=/home/ubuntu/recostock/deploy/.env
ExecStart=/home/ubuntu/recostock/venv/bin/python scripts/run_intraday.py
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=intraday-bot

[Install]
WantedBy=multi-user.target
